# soundboard.service
[Unit]
Description=Sound board service

[Service]
ExecStart=/home/pi/Soundboard/soundbox.sh
[Install]
WantedBy=multi-user.target
